# services/node-exporter/run
#!/bin/zsh
set -euo pipefail

exec node_exporter 2>&1
